# modern-kiosk-user.service
[Unit]
Description=Kimonokittens Kiosk Browser (User Service)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
Environment="XDG_RUNTIME_DIR=/run/user/1001"
ExecStartPre=/bin/sleep 15
ExecStart=/usr/bin/google-chrome --kiosk --no-first-run --disable-infobars --disable-session-crashed-bubble --disable-web-security --disable-features=TranslateUI --noerrdialogs --incognito --no-default-browser-check --password-store=basic --start-maximized --app=http://localhost
Restart=always
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=300

[Install]
WantedBy=default.target